# deployment/spacetimedb.service
# SpacetimeDB systemd user service
# Install to: ~/.config/systemd/user/spacetimedb.service
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp deployment/spacetimedb.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable spacetimedb.service
#   systemctl --user start spacetimedb.service
#
# Check status:
#   systemctl --user status spacetimedb.service
#   journalctl --user -u spacetimedb.service -f

[Unit]
Description=SpacetimeDB Server for Sunaba Multiplayer
Documentation=https://spacetimedb.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Restart on failure
Restart=on-failure
RestartSec=10s

# Working directory for SpacetimeDB data
WorkingDirectory=%h/spacetimedb
# Create data directory if it doesn't exist (run manually first time)
# mkdir -p ~/spacetimedb

# Admin whitelist (comma-separated email list for admin access)
# Uncomment and set to enable admin-only features
# Environment="SUNABA_ADMIN_EMAILS=your.email@gmail.com,friend@gmail.com"

# Start SpacetimeDB listening on all interfaces
# Port 3000 for HTTP + WebSocket connections
ExecStart=/usr/local/bin/spacetime start --listen-addr=0.0.0.0:3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spacetimedb

# Security hardening (optional but recommended)
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=default.target
